# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=gunicorn -k eventlet chat:app -b 0.0.0.0:80
stdout_logfile=/flask.log
stderr_logfile=/flask.log
autostart=true
autorestart=true